# ops/vps/nexova-tool-watch.service
[Unit]
Description=Restart nexova Hermes gateway when the agent tool manifest changes
# Don't run while the gateway itself is starting up.
After=network-online.target

[Service]
Type=oneshot
# Invoke via bash explicitly rather than relying on the script's shebang — a
# heredoc paste can drop the `#!` line, which makes systemd fail with
# 203/EXEC "Exec format error" even though the script body is fine.
ExecStart=/bin/bash %h/.hermes/profiles/nexova/nexova-tool-watch.sh
